# deploy/drop.service
[Unit]
Description=drop personal capture API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/drop
EnvironmentFile=/opt/drop/.env
ExecStart=/opt/drop/.venv/bin/uvicorn drop.app:app --host 127.0.0.1 --port 9731
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
